# contrib/init/motiond.conf
description "Motion Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env MOTIOND_BIN="/usr/bin/motiond"
env MOTIOND_USER="motioncore"
env MOTIOND_GROUP="motioncore"
env MOTIOND_PIDDIR="/var/run/motiond"
# upstart can't handle variables constructed with other variables
env MOTIOND_PIDFILE="/var/run/motiond/motiond.pid"
env MOTIOND_CONFIGFILE="/etc/motioncore/motion.conf"
env MOTIOND_DATADIR="/var/lib/motiond"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # motiond will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$MOTIOND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run motiond."
        echo "The setting must appear in $MOTIOND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Motion Core Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$MOTIOND_PIDDIR"
    chmod 0755 "$MOTIOND_PIDDIR"
    chown $MOTIOND_USER:$MOTIOND_GROUP "$MOTIOND_PIDDIR"
    chown $MOTIOND_USER:$MOTIOND_GROUP "$MOTIOND_CONFIGFILE"
    chmod 0660 "$MOTIOND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$MOTIOND_PIDFILE" \
    --chuid $MOTIOND_USER:$MOTIOND_GROUP \
    --exec "$MOTIOND_BIN" \
    -- \
    -pid="$MOTIOND_PIDFILE" \
    -conf="$MOTIOND_CONFIGFILE" \
    -datadir="$MOTIOND_DATADIR" \
    -disablewallet \
    -daemon

